# meta-resticsrv/recipes-core/base-files/base-files/ledsoff.rules
SUBSYSTEM=="leds", ACTION=="add", RUN+="/bin/sh -c 'echo 0 > /sys%p/brightness'"
